# dnf-autoclean.service
[Unit]
Description=DNF AutoClean Service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/dnf-autoclean.py
