# dnf-autoclean.service
[Unit]
Description=DNF AutoClean service

[Service]
Type=oneshot
ExecStart=/usr/bin/dnf-autoclean
